# deploy/eqrf.service.example
[Unit]
Description=EQRF local network application
After=network-online.target
Wants=network-online.target

[Service]
User=eqrf
Group=eqrf
WorkingDirectory=/opt/EQRF
EnvironmentFile=/opt/EQRF/.env
ExecStart=/opt/EQRF/venv/bin/gunicorn --worker-class gthread --workers 1 --threads 24 --timeout 0 --bind 0.0.0.0:8000 wsgi:application
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/EQRF/data /opt/EQRF/pdfs /opt/EQRF/static /opt/EQRF/backups
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
