# deploy/eqrf-backup.timer.example
[Unit]
Description=Run EQRF backup daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=eqrf-backup.service

[Install]
WantedBy=timers.target
